# docker/postgres/cron.ini
[supervisord]
logfile=/var/log/supervisord.log
loglevel=info
nodaemon=true

[program:postgres]
command=/usr/local/bin/docker-entrypoint.sh postgres -c config_file=/etc/postgresql.conf
autostart=true
autorestart=true

[program:cron]
command=/usr/sbin/crond -f -l 8
autostart=true
autorestart=true