# docker/postgres/crontab
*/5 * * * * remove_tests > /var/log/remove_tests.log